# rootfs/filesystem/etc/sv/tty2/run
#!/bin/sh
exec busybox getty -L tty2 115200 vt100 -l /bin/login
