# rootfs/filesystem/etc/sv/tty7/run
#!/bin/sh
exec busybox getty -L tty7 115200 vt100 -l /bin/login
