# rootfs/filesystem/etc/sv/mdev/run
#!/bin/sh

busybox mdev -s -S -f -d 2>&1
exec sleep infinity
