# rootfs/filesystem/etc/sv/syslogd/run
#!/bin/sh

mkdir -p /var/log
exec busybox syslogd -n -O /var/log/messages 2>&1
