# rootfs/filesystem/etc/sv/tty1/run
#!/bin/sh

sleep 0.1

exec busybox getty -L tty1 115200 vt100 -l /bin/login
