# rootfs/filesystem/etc/service/mdev-service/run
#!/bin/sh

printf "\033[96m*\033[0m Starting \033[97mmdev-service\033[0m\n"

echo /sbin/mdev > /proc/sys/kernel/hotplug
busybox mdev -s
touch /etc/service/mdev-service/.lock

exec sleep infinity
