# rootfs/filesystem/etc/sv/hostname-reload-daemon/run
#!/bin/sh

exec /sbin/reload-hostname-daemon
