# rootfs/filesystem/etc/service/hostname-reload-daemon/run
#!/bin/sh

printf "\033[96m*\033[0m Starting \033[97mhostname-reload-daemon\033[0m\n" > /dev/console
exec /sbin/reload-hostname-daemon
