# rootfs/filesystem/etc/service/tty6/run
#!/bin/sh
exec busybox getty -L tty6 115200 vt100 -l /bin/login
