# rootfs/filesystem/etc/service/klogd-service/run
#!/bin/sh
touch /etc/service/klogd-service/.lock
exec busybox klogd -n
