# rootfs/filesystem/etc/sv/tty4/run
#!/bin/sh
exec busybox getty -L tty4 115200 vt100 -l /bin/login
