# rootfs/filesystem/etc/service/klogd/run
#!/bin/sh
printf "\033[96m*\033[0m Starting \033[97mklogd\033[0m\n"

exec busybox klogd -n
