# rootfs/filesystem/etc/sv/tty5/run
#!/bin/sh
exec busybox getty -L tty5 115200 vt100 -l /bin/login
