# rootfs/filesystem/etc/service/dhcp/run
#!/bin/sh

printf "\033[96m*\033[0m Starting \033[97mdhcp\033[0m\n" > /dev/console

ip link set eth0 up
exec 2>&1
exec udhcpc -i eth0 -f -S
