# rootfs/filesystem/etc/sv/dhcp/run
#!/bin/sh

printf "\033[96m*\033[0m Starting \033[97mdhcp\033[0m\n" > /dev/console

ip link set eth0 up
mkdir -p /var/log/dhcp
exec udhcpc -i eth0 -f -S 2>&1
