# rootfs/filesystem/etc/sv/tty3/run
#!/bin/sh
exec busybox getty -L tty3 115200 vt100 -l /bin/login
